# logs/poem-automation.service
[Unit]
Description=Daily Poem Generation Service
After=network.target

[Service]
Type=simple
User=barroca888
WorkingDirectory=/Users/barroca888/Downloads/Dev/Git Ai Poem Agent/Git-AI-Poem-Agent
ExecStart=/opt/anaconda3/bin/python daily_automation.py
Restart=always
RestartSec=3600
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 